# systemd/yt3-discord.service
[Unit]
Description=YT3 Discord Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/kafka/2511youtuber/v3/yt3
ExecStart=/usr/local/bin/bun src/agents/standalone/discord_bot.ts
Restart=on-failure
RestartSec=5
Environment=PATH=/home/kafka/.bun/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/kafka

[Install]
WantedBy=default.target
